# setup /toypad-usb-gadget.service
[Unit]
Description=Toy Pad USB Gadget Setup
Before=toypad-emulator.service
After=sys-kernel-config.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/toypad-usb-gadget-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=toypad-usb-gadget

[Install]
WantedBy=multi-user.target
