# setup/toypad-emulator.service
[Unit]
Description=Multi-Platform Toy Pad Emulator
Documentation=https://github.com/yourusername/multi_toy_pad_emulator
After=network.target
Wants=toypad-usb-gadget.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/multi_toy_pad_emulator
Environment="NODE_ENV=production"
ExecStart=/usr/local/bin/node /home/pi/multi_toy_pad_emulator/dist/server/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=toypad-emulator

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
